# romm/rootfs/etc/s6-overlay/s6-rc.d/rq-scheduler/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting RQ scheduler..."

cd /backend || bashio::exit.nok "Failed to change to /backend directory"

exec rqscheduler \
    --url "${REDIS_URL}" \
    --pid /tmp/rq_scheduler.pid \
    --interval 60
