# profilarr-v2/rootfs/etc/services.d/profilarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Profilarr V2 service
# ==============================================================================

exec 2>&1

bashio::log.info "Starting Profilarr V2..."

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Set BASE_URL/SCRIPT_NAME for ingress support
# The Profilarr binary reads environment variables for routing
if bashio::addon.ingress; then
    export BASE_URL="$(bashio::addon.ingress_entry)"
    export SCRIPT_NAME="$(bashio::addon.ingress_entry)"
    bashio::log.info "Ingress enabled, BASE_URL set to: ${BASE_URL}"
fi

# Start Profilarr V2 binary directly (no wrapper needed)
# The binary is a compiled Deno application that handles its own user/group via PUID/PGID
exec /app/profilarr
