# .claude/skills/ha-apps/scaffold/v3_example/s6-overlay/s6-rc.d/example-init/up
/etc/s6-overlay/scripts/example-init.sh
