# romm/rootfs/etc/s6-overlay/s6-rc.d/rq-worker-log/run
#!/bin/sh
exec logutil-service /var/log/rq-worker
