# romm/rootfs/etc/s6-overlay/s6-rc.d/valkey-log-prepare/up
if { mkdir -p /var/log/valkey }
chmod 02755 /var/log/valkey
