# profilarr/rootfs/etc/services.d/profilarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Profilarr service
# ==============================================================================

bashio::log.info "Starting Profilarr..."

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Set BASE_URL/SCRIPT_NAME for ingress support
# Get the ingress entry path so the Flask app knows the base URL for routing
if bashio::addon.ingress; then
    export BASE_URL="$(bashio::addon.ingress_entry)"
    export SCRIPT_NAME="$(bashio::addon.ingress_entry)"
    bashio::log.info "Ingress enabled, BASE_URL/SCRIPT_NAME set to: ${BASE_URL}"
fi

# Start Profilarr directly with gunicorn (skip upstream entrypoint)
exec gunicorn \
    --bind 0.0.0.0:6868 \
    --timeout 600 \
    --workers 1 \
    --access-logfile - \
    --error-logfile - \
    "app.main:create_app()"
