# profilarr/rootfs/etc/services.d/profilarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Profilarr service
# ==============================================================================

bashio::log.info "Starting Profilarr..."

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Set SCRIPT_NAME for ingress support
# This tells Flask/WSGI the base path for routing when accessed through HA ingress
if bashio::addon.ingress; then
    export SCRIPT_NAME="$(bashio::addon.ingress_entry)"
    bashio::log.info "Ingress enabled, SCRIPT_NAME set to: ${SCRIPT_NAME}"
fi

# Start Profilarr directly with gunicorn (skip upstream entrypoint)
exec gunicorn \
    --bind 0.0.0.0:6868 \
    --timeout 600 \
    --workers 1 \
    --access-logfile - \
    --error-logfile - \
    "app.main:create_app()"
