# profilarr/rootfs/etc/services.d/profilarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Profilarr service
# ==============================================================================

bashio::log.info "Starting Profilarr..."

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Get mapped port from Home Assistant
PORT="$(bashio::addon.port 6868)"

# Start Profilarr directly with gunicorn (skip upstream entrypoint)
exec gunicorn \
    --bind 0.0.0.0:"${PORT}" \
    --timeout 600 \
    --workers 1 \
    --access-logfile - \
    --error-logfile - \
    "app.main:create_app()"
