# profilarr/rootfs/etc/services.d/profilarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Profilarr service
# ==============================================================================

bashio::log.info "Starting Profilarr..."

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Always bind to localhost:6868 for nginx proxy
# Nginx handles ingress access on the configured ingress_port
PORT=6868

bashio::log.info "Profilarr will listen on 0.0.0..${PORT}"

# Start Profilarr directly with gunicorn (skip upstream entrypoint)
exec gunicorn \
    --bind 0.0.0.0:"${PORT}" \
    --timeout 600 \
    --workers 1 \
    --access-logfile - \
    --error-logfile - \
    "app.main:create_app()"
