# romm/rootfs/etc/s6-overlay/s6-rc.d/romm-main-log-prepare/up
if { mkdir -p /var/log/romm-main }
chmod 02755 /var/log/romm-main
