# romm/rootfs/etc/s6-overlay/s6-rc.d/rq-scheduler-log-prepare/up
if { mkdir -p /var/log/rq-scheduler }
chmod 02755 /var/log/rq-scheduler
