# .claude/skills/ha-apps/templates/ingress-nginx-s6-v3/rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Application service finish script
# s6-overlay v3: This script runs when your application exits
# Arguments: $1 = exit code (256 = killed by signal)
# ==============================================================================

set -e

declare APP_EXIT_CODE=${1}

bashio::log.info "Application service exited with code ${APP_EXIT_CODE}"

# Exit code meanings:
# 0    = Normal shutdown (user requested)
# 256  = Killed by signal (SIGTERM during container shutdown)
# Other = Abnormal exit (error/crash)

if [[ "${APP_EXIT_CODE}" -eq 0 ]]; then
    bashio::log.info "Application shut down normally"
    # Let container stop normally
elif [[ "${APP_EXIT_CODE}" -eq 256 ]]; then
    bashio::log.info "Application was terminated by signal"
    # Let container stop normally
else
    # Application crashed - halt the container
    bashio::log.error "Application crashed with exit code ${APP_EXIT_CODE}"
    bashio::log.error "Halting container to prevent restart loop"

    # Write exit code for container to return
    echo "${APP_EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode

    # Halt the entire s6 supervision tree
    # This stops all services and exits the container
    exec /run/s6/basedir/bin/halt
fi
