# romm/rootfs/etc/s6-overlay/s6-rc.d/nginx-log-prepare/up
if { mkdir -p /var/log/nginx-service }
chmod 02755 /var/log/nginx-service
