# romm/rootfs/etc/services.d/valkey/finish
#!/usr/bin/with-contenv bashio
declare APP_EXIT_CODE=${1}

bashio::log.info "Stopping valkey..."

if [ -f /tmp/valkey-server.pid ]; then
    kill $(cat /tmp/valkey-server.pid) 2>/dev/null || true
    rm -f /tmp/valkey-server.pid
fi
