# romm/rootfs/etc/s6-overlay/s6-rc.d/rq-worker-log-prepare/up
if { mkdir -p /var/log/rq-worker }
chmod 02755 /var/log/rq-worker
