# romm/rootfs/etc/s6-overlay/s6-rc.d/rq-worker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting RQ worker..."

cd /backend || bashio::exit.nok "Failed to change to /backend directory"

exec rq worker \
    --path /backend \
    --url "${REDIS_URL}" \
    --results-ttl 86400 \
    --pid /tmp/rq_worker.pid \
    high default low
