# romm/rootfs/etc/s6-overlay/s6-rc.d/valkey/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting internal valkey..."

exec valkey-server \
    --dir /data/redis_data \
    --daemonize no \
    --loglevel warning \
    --pidfile /tmp/valkey-server.pid
