# template/rootfs/etc/services.d/myapp/run
#!/usr/bin/with-contenv bashio
# =============================================================================
# Home Assistant Add-on: Your Add-on
# Main service script - runs the application
# =============================================================================

bashio::log.info "Starting Your Add-on..."

# Set up error handling
set -e

# Change to application directory (if applicable)
# cd /app || bashio::exit.nok "Failed to change to application directory"

# =============================================================================
# Pre-start checks
# =============================================================================

# Example: Wait for dependent service
# bashio::net.wait_for 4506 30 2>&1 || bashio::exit.nok "Database not available"

# Example: Check required files exist
# if [ ! -f "/data/config.yml" ]; then
#     bashio::log.error "Configuration file not found!"
#     bashio::exit.nok "Missing configuration"
# fi

# =============================================================================
# Start your application
# =============================================================================

# Example 1: Python application
# exec python3 /app/main.py

# Example 2: Node.js application
# exec node /app/index.js

# Example 3: Background service with logging
# exec your_application \
#     --log-level="${LOG_LEVEL}" \
#     --port="${APP_PORT}" \
#     2>&1

# Example 4: Simple long-running script
# while true; do
#     bashio::log.info "Application running..."
#     sleep 60
# done

# Placeholder - replace with your actual application start command
bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Application started successfully"

# Keep container running (replace with actual application)
exec sleep infinity
